# diary-backup.service
[Unit]
Description=Diary database backup

[Service]
Type=oneshot
ExecStart=/opt/diary/backup.sh
